# barrier_breaker/target/linux/mvebu/mamba/base-files/etc/init.d/fan
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org 

START=99
start() {
	# Run leds_monitor
	/usr/sbin/sensors_monitor&
}
