# barrier_breaker/package/linksys-base-files/files/etc/hotplug.d/block/11-smb-share
device=`basename $DEVPATH`
{ echo ${device} | grep -q mtd; } || \
	/sbin/smb_share.sh "${ACTION}" /dev/${device}
